# lcd.service
[Unit]
Description=lcd operation daemon

[Service]
Type=simple
ExecStart=/usr/local/bin/lcd-daemon.py

[Install]
WantedBy=multi-user.target
